# systemd/elasticsearch.service
[Unit]
Description=Elasticsearch for Awful Security News
Requires=docker.service
After=docker.service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/tg/awful_security_news

# Start the container
ExecStart=/usr/bin/docker compose -f /home/tg/awful_security_news/docker/docker-compose.prod.yml up -d

# Stop the container
ExecStop=/usr/bin/docker compose -f /home/tg/awful_security_news/docker/docker-compose.prod.yml down

# Restart policy
Restart=on-failure
RestartSec=10s

# Security settings
User=tg
Group=tg

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
